# roles/kde/files/devilspie2.service
[Unit]
Description=Devilspie

[Service]
ExecStart=/usr/bin/devilspie2

[Install]
WantedBy=default.target
