# roles/kde/files/devilspie2.service
[Unit]
Description=Devilspie2

[Service]
ExecStart=/usr/bin/devilspie2
Restart=always

[Install]
WantedBy=default.target
